# etc/systemd/system/recovery-check.service
[Unit]
Description=Recovery Script Runner
Before=emulationstation.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/recovery-runner.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target